# deploy/ce-worker.service
[Unit]
Description=CE native worker (shares cores with the mesh via the local ce node)
After=network-online.target ce-relay.service
Wants=network-online.target

[Service]
# Attaches to the ce node running on this same host (its HTTP+SSE API) and joins the
# mesh as the "ce-worker" service; jobs arrive over libp2p, not a hub WebSocket.
ExecStart=/usr/bin/node /opt/ce-worker/worker.js --node http://127.0.0.1:8844 --name relay
WorkingDirectory=/opt/ce-worker
Restart=always
RestartSec=5
# Keep the relay responsive: the worker shares cores but never starves bootstrap/relay duty.
Nice=10
CPUWeight=50

[Install]
WantedBy=multi-user.target
