# deploy/com.cenet.worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd agent: run a CE native worker in the background, no browser.
  Install:
    cp deploy/com.cenet.worker.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.cenet.worker.plist
  Stop:
    launchctl unload -w ~/Library/LaunchAgents/com.cenet.worker.plist
  Edit the path below if this repo is not at ~/ce-net.
-->
<plist version="1.0">
<dict>
  <key>Label</key>            <string>com.cenet.worker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/07lead01/.nvm/versions/node/v22.12.0/bin/node</string>
    <string>/Users/07lead01/ce-net/ce-worker/worker.js</string>
    <string>--hub</string>
    <string>wss://ce-net.com/hub</string>
  </array>
  <key>RunAtLoad</key>        <true/>
  <key>KeepAlive</key>        <true/>
  <key>ProcessType</key>      <string>Background</string>
  <key>StandardOutPath</key>  <string>/tmp/ce-worker.log</string>
  <key>StandardErrorPath</key><string>/tmp/ce-worker.err</string>
</dict>
</plist>
